# bin/assets/app.conf
[program:app]
directory=/var/www/zero-to-sixty
command=npm start
stdout_logfile=/var/log/zts.out
stderr_logfile=/var/log/zts.err
stopasgroup=true
stopsignal=QUIT
