# bin/assets/app.conf
[program:app]
directory=/var/www/zero-to-sixty
command=npm start
